# units/yourls@.service.d/https-site.conf
# Use the HTTPS site instead of http
[Service]
# Set the site based on the instance name
Environment='DOCKER_ENV_YOURLS_SITE=YOURLS_SITE=https://%i'
